# photoframegtk.service
[Unit]
Description=Photo Frame GTK service file
After=graphical.target

[Service]
ExecStart=DISPLAY=:0 /usr/bin/python3 /home/hpi/photo_frame_gtk/photoframegtk.py
User=root
Group=root


[Install]
WantedBy=multi-user.target
